# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:qdrant]
command=/usr/local/bin/qdrant
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:api]
command=gunicorn -k uvicorn.workers.UvicornWorker main:app --bind 0.0.0.0:10000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
